# integration/benchmark/server/systemd/flow-tps.service
[Unit]
Description=Flow TPS tests - generate list of merge commit hashes and run TPS tests against each one
After=network.target

[Service]
Type=simple
ExecStart=/var/flow/flow-go/integration/benchmark/server/tps.sh
WorkingDirectory=/var/flow/flow-go/integration/benchmark/server
Environment="GOPATH=/var/flow/go" "GOCACHE=/var/flow/gocache"
Restart=always

[Install]
WantedBy=multi-user.target
