# integration/benchmark/server/systemd/flow-tps.timer
[Unit]
Description=Run Flow TPS tests once per day

[Timer]
OnUnitActiveSec=1440min

[Install]
WantedBy=timers.target
