# integration/benchmark/server/systemd/flow-tps.timer
[Unit]
Description=Run control.sh every day

[Timer]
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
