# integration/benchmark/server/systemd/flow-tps.service
[Unit]
Description=Flow TPS tests - generate list of merge commit hashes and run TPS tests against each one

[Service]
Type=oneshot
ExecStart=/opt/flow-go/integration/benchmark/server/control.sh
ExecStart=/opt/flow-go/integration/benchmark/server/hello.sh
ExecStart=/opt/flow-go/integration/benchmark/server/bench.sh
WorkingDirectory=/opt/flow-go/integration/benchmark/server
Environment="GOPATH=/opt/go" "GOCACHE=/opt/gocache"
RemainAfterExit=no
