# integration/benchmark/server/systemd/flow-tps.service
[Unit]
Description=Flow TPS tests - generate list of merge commit hashes and run TPS tests against each one

[Service]
Type=oneshot
ExecStart=/mnt/sdb/flow-go/integration/benchmark/server/runs.sh
ExecStart=/mnt/sdb/flow-go/integration/benchmark/server/control.sh
ExecStart=/mnt/sdb/flow-go/integration/benchmark/server/bench.sh
WorkingDirectory=/mnt/sdb/flow-go/integration/benchmark/server
Environment="GOPATH=/mnt/sdb/go" "GOCACHE=/mnt/sdb/gocache"
RemainAfterExit=no
