# integration/benchmark/server/systemd/flow-tps.service
[Unit]
Description=Flow TPS tests - control script to generate list of merge hashes

[Service]
Type=oneshot
ExecStart=/home/user/flow-go/integration/benchmark/server/control.sh
WorkingDirectory=/tmp
RemainAfterExit=no
